# vps-bot.service
[Unit]
Description=DarkNodes VPS Manager Discord Bot
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root

# Switched to the global system Python interpreter
WorkingDirectory=/root/Mybot
ExecStart=/usr/bin/python3 bot.py

# Auto-recovery configuration
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=vps-bot

[Install]
WantedBy=multi-user.target
