# vps-bot.service
[Unit]
Description=DarkNodes VPS Manager Discord Bot
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
Group=root

# Dynamically mapped repository paths (resolved during installation)
WorkingDirectory=/root/Mybot
ExecStart=/root/Mybot/.venv/bin/python3 bot.py

# Auto-recovery if the bot hits a rare crash or memory limit
Restart=always
RestartSec=5

# Console logging integration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vps-bot

[Install]
WantedBy=multi-user.target
